# Library/LaunchAgents/com.yarden.fix-chrome-password-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yarden.fix-chrome-password-manager</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/kcw/bin/fix-chrome-password-manager.sh</string>
  </array>

  <key>LimitLoadToSessionType</key>
  <array>
    <string>Aqua</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/kcw/Library/Logs/fix-chrome-password-manager.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kcw/Library/Logs/fix-chrome-password-manager.log</string>
</dict>
</plist>
